# deploy/systemd/settle-listener.service
[Unit]
Description=Settlepay blockchain listener (eth_getLogs polling per chain)
Documentation=https://github.com/kindrakevich-agency/settlepay.pro
After=network.target mariadb.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/www/wwwroot/settlepay.pro
Environment=APP_ENV=prod

# Mainnet — watches Base, Polygon, Arbitrum, Optimism. For dev / dogfood
# testing on sepolia testnets, enable settle-listener-testnet.service
# alongside this one (they don't conflict — different chain_ids → different
# chain_cursors rows in the DB).
ExecStart=/www/server/php/83/bin/php bin/console app:chain:listen

Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

StandardOutput=append:/var/log/settle/listener.log
StandardError=append:/var/log/settle/listener-error.log

# Resource caps so a runaway listener can't take the box down.
MemoryLimit=256M

[Install]
WantedBy=multi-user.target
